# hytte.service
[Unit]
Description=Hytte web server
After=network.target

[Service]
Type=simple
User=robin
Group=robin
WorkingDirectory=/home/robin/Hytte
ExecStart=/home/robin/Hytte/bin/hytte
Restart=always
RestartSec=5
Environment=PORT=8080
EnvironmentFile=-/home/robin/Hytte/.env

[Install]
WantedBy=multi-user.target
